# Library/LaunchAgents/sh.brew.gnupg.gpg-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<!-- man launchd -->
	<!-- Set the name of the Launchd task -->
	<key>Label</key>
	<string>sh.brew.gnupg.gpg-agent</string>
	<!-- Set a command to run and its arguments -->
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/gpgconf</string>
		<string>--launch</string>
		<string>gpg-agent</string>
	</array>
	<!-- Tells the task to run once the XML is loaded -->
	<key>RunAtLoad</key>
	<true />
	<!-- <key>Sockets</key>
	<dict>
		<key>Listeners</key>
		<dict>
			<key>SecureSocketWithKey</key>
			<string>SSH_AUTH_SOCK</string>
		</dict>
	</dict>
	<key>EnableTransactions</key>
	<true/> -->
</dict>
</plist>
